# Systems/03pillarsSoundFX/pillar-audio.service
[Unit]
Description=Pharaohs Pillar Serial Audio Service
After=network-online.target sound.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/pi
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=-/home/pi/pillar_audio.env
ExecStart=/usr/bin/python3 /home/pi/pillar_audio_serial.py --port /dev/ttyS0 --baud 115200 --device plughw:CARD=Device,DEV=0
Restart=always
RestartSec=2
User=pi
Group=audio

[Install]
WantedBy=multi-user.target
